# terraform/30_vms/vm/ignition/mounts/contents/systemd/tpl.mount
[Mount]
What=${source}
Where=${target}
Type=${type}

[Install]
WantedBy=local-fs.target